# overlays/etc/systemd/system/setup.service
[Unit]
Description=First Boot Setup Service
After=network.target
Before=graphical.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/first-boot-setup.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
