# overlays/etc/systemd/system/wifi-tools.service
[Unit]
Description=WiFi Tools Background Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/wifi-tools-daemon
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
